# ansible/roles/flask/templates/flask_app.service.j2
[Unit]
Description=Flask app

[Service]
User=nginx
WorkingDirectory={{app_dir}}
Environment="PATH=/usr/bin"
ExecStart=/usr/bin/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target